# rc.local
#!/bin/bash
# THIS FILE IS ADDED FOR COMPATIBILITY PURPOSES
#
# It is highly advisable to create own systemd services or udev rules
# to run scripts during boot instead of using this file.
#
# In contrast to previous versions due to parallel execution during boot
# this script will NOT be run after all other services.
#
# Please note that you must run 'chmod +x /etc/rc.d/rc.local' to ensure
# that this script will be executed during boot.

su - poker -c "/home/poker/goes17_ralph/process_ABI_rgb_realtime-devel-python3.6_conus_meso/run_rgb_fulldisk_latest.sh >& /home/poker/goes17_ralph/process_ABI_rgb_realtime-devel-python3.6_conus_meso/run_rgb_fulldisk_latest.out& "
su - poker -c "/home/poker/goes17/fulldiskc_grb/run_ir13_fulldisk_loop.sh >& /home/poker/goes17/fulldiskc_grb/run_ir13_fulldisk_loop.out & "
su - poker -c "/home/poker/goes17/fulldiskc_grb/run_ir14_fulldisk_loop.sh >& /home/poker/goes17/fulldiskc_grb/run_ir14_fulldisk_loop.out & "
su - poker -c "/home/poker/goes17/fulldiskc_grb/run_wvc_fulldisk_loop.sh >& /home/poker/goes17/fulldiskc_grb/run_wvc_fulldisk_loop.out & "
su - poker -c "/home/poker/goes17/fulldiskc_grb/run_wvh_fulldisk_loop.sh >& /home/poker/goes17/fulldiskc_grb/run_wvh_fulldisk_loop.out & "
su - poker -c "/home/poker/goes17/fulldiskc_grb/run_wvl_fulldisk_loop.sh >& /home/poker/goes17/fulldiskc_grb/run_wvl_fulldisk_loop.out & "

touch /var/lock/subsys/local
